# admin/thumbsupdbwrite.timer
[Unit]
Description=Run Thumbsup database write every minute

[Timer]
OnBootSec=15s
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
