# admin/thumbsuptoken.timer
[Unit]
Description=Run Thumbsup token generator every 55 minutes

[Timer]
OnBootSec=10s
OnUnitActiveSec=55m

[Install]
WantedBy=timers.target
